# etc/systemd/system/addWebsites.path
[Path]
PathModified=/etc/websites

[Install]
WantedBy=multi-user.target
